# tproxyrule-systemd/tproxyrule.service
[Unit]
Description=TPROXY Rule
After=network.target network-online.target

[Service]
ExecStart=/usr/lib/tproxy-rule
Type=oneshot

[Install]
WantedBy=multi-user.target
